# crontab
* * * * * /bin/bash /home/embassy/schedule.sh >> /home/embassy/logfile.log
