# crontab/main
* * * * * /bin/bash /home/embassy/shell_scripts/schedule.sh >> /home/embassy/logfile.log
